# conf/sdpmc.conf
# sdpmc - SDP Master Controller
#
# The SDP master controller is the first SDP service to be started and coordinates
# all other SDP processes.

description	"SDP Master Controller"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/bin/sdp_master_controller.py
