# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:sdpmc]
command=/usr/local/bin/sdp_master_controller.py

[program:redis]
command=/usr/bin/redis-server
